# scripts/sv/API/run
#!/bin/bash -e
exec 2>&1
#API_DIR=harvester_gRPC
exec chpst -u ads bash -c 'cd /app/$APP_DIR/ && python3 run.py HARVESTER_API'