# scripts/sv/API/run
#!/bin/bash -e
exec 2>&1
exec chpst -u ads bash -c 'cd /app/SciXHarvester && python3 run.py HARVESTER_API'