# scripts/sv/APP/run
#!/bin/bash -e
exec 2>&1
#APP_DIR=harvester_gRPC
exec chpst -u ads bash -c 'cd /app/$APP_DIR/ && alembic upgrade head && python3 /app/scripts/kafka_create_topics.py  && python3 run.py HARVESTER_APP'