# scripts/sv/APP/run
#!/bin/bash -e
exec 2>&1
#APP_DIR=harvester_gRPC
exec chpst -u ads bash -c 'cd /app/$APP_DIR/ && alembic upgrade head && python3 run.py'